# package/DEBIAN/prerm
#!/bin/sh

systemctl stop fan-control
systemctl disable fan-control
